# updater/updater.service
[Unit]
Description=Script to update dotfiles

[Service]
ExecStart=/bin/dotfile-updater
